# vagrant/ansible/roles/apache/files/apache.cron
# Ansible inserted crontab for graceful restarting Apache

PATH=/usr/sbin:/usr/bin:/sbin:/bin

@reboot      root service apache2 graceful >> /var/log/cron.log 2>&1
*/30 * * * * root service apache2 graceful >> /var/log/cron.log 2>&1

